# lighttpd-unicorn/lighttpd-unicorn.service
[Unit]
Description=A secure, fast, compliant and very flexible web-server (unicorn proxy)
After=syslog.target network.target

[Service]
PrivateTmp=true
ExecStart=/usr/sbin/lighttpd-angel -D -f /etc/lighttpd/lighttpd-unicorn.conf
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
